# homestation.service
[Unit]
Description=HomeStation Pi4 Hub
After=network-online.target mosquitto.service
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/homestation
ExecStart=/usr/bin/python3 /home/pi/homestation/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
